# services/moder/moder.service
[Unit]
Description=DImROD Service - Moder
Documentation=https://github.com/cwshugg/dimrod
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/path/to/moder.py --config /path/to/config.json --oracle

[Install]
WantedBy=multi-user.target

