# cron/crontab
15 12 * * * /home/provolone/chs/cron/status-email.sh

