# cron/crontab
15 12 * * * provolone /home/provolone/chs/status-email.sh

